# deployment/dxid-validator.service
[Unit]
Description=dxID Layer0 Validator Node
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/dxid-layer0
ExecStart=/home/ubuntu/dxid-layer0/target/release/dxid-node \
    --p2p-listen 0.0.0.0:7000 \
    --rpc-listen 0.0.0.0:8545 \
    --validator \
    --bootstrap-peer YOUR_VM2_IP:7000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dxid-validator

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/dxid-layer0/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
